# Resources/com.maatwork.meetcapture.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maatwork.meetcapture.daemon</string>

    <!-- Phase 6 fix: use bundled meet-daemon shell wrapper (real path).
         The previous BundleProgram pointing to daemon_main.py was a phantom
         reference; daemon_main.py was deleted in the 2026-06-04 refactor.
         The shell wrapper resolves the active Python (bundled venv → homebrew → system)
         and execs server.py inside Contents/Resources. -->
    <key>ProgramArguments</key>
    <array>
        <string>__BUNDLE_DIR__/Contents/Resources/meet-daemon</string>
    </array>

    <!-- Run at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on crash, not on clean exit -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Prevent crash-loop restarts -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Background priority -->
    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
    <key>Nice</key>
    <integer>10</integer>

    <!-- Phase 6: clean shutdown grace period -->
    <key>ExitTimeOut</key>
    <integer>30</integer>

    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/meetcapture-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/meetcapture-daemon.log</string>
</dict>
</plist>
